# deploy/systemd/secretary.service
[Unit]
Description=Secretary MVP
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/secretary
EnvironmentFile=/root/secretary/.env
ExecStart=/root/secretary/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target